# supervisor.conf
[program:<project_name>]
environment=DJANGO_SETTINGS_MODULE='core.settings.production'
command=/bin/gunicorn_start <project_name>
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/<project_name>_stdout.log
stderr_logfile=/var/log/supervisor/<project_name>_stderr.log
